# recipes/weston/files/weston.socket
[Socket]
ListenStream=/tmp/wayland-early

[Install]
WantedBy=sockets.target

