# recipes/weston/files/weston.socket
[Socket]
ListenStream=/tmp/wayland-secondary

[Install]
WantedBy=sockets.target
